# files/etc/uci-defaults/99-aq-init
#!/bin/sh

/etc/init.d/aq-init enabled
/etc/init.d/aq-init start

exit 0